# Oct21/chef-starter/chef-repo/cookbooks/tomcat9/templates/tomcat.service.erb
[Unit]
Description=Tomcat 9 servlet container
After=network.target

[Service]
Type=forking

User=<%= node['tomcat9']['username'] %>
Group=<%= node['tomcat9']['groupname'] %>

Environment="JAVA_HOME=<%= node['tomcat9']['java_home'] %>"
Environment="JAVA_OPTS=<%= node['tomcat9']['java_opts'] %>"

Environment="CATALINA_BASE=<%= node['tomcat9']['userhome'] %>/latest"
Environment="CATALINA_HOME=<%= node['tomcat9']['userhome'] %>/latest"
Environment="CATALINA_PID=<%= node['tomcat9']['userhome'] %>/latest/temp/tomcat.pid"
Environment="CATALINA_OPTS=-Xms512M -Xmx1024M -server -XX:+UseParallelGC"

ExecStart=<%= node['tomcat9']['userhome'] %>/latest/bin/startup.sh
ExecStop=<%= node['tomcat9']['userhome'] %>/latest/bin/shutdown.sh

[Install]
WantedBy=multi-user.target